# config/includes.chroot_after_packages/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin sinco --noclear %I $TERM
